# docker/php/crontab
* * * * * root echo "Hello Romario" >> /var/log/cron.log 2>&1
